# file/rc.local
#!/bin/sh -e
ssserver -c /home/wangwaking/ss.json -d start&
